# deploy/sportuniverse-queue.service
[Unit]
Description=SportUniverse Laravel queue worker
After=network.target redis-server.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/sportuniverse
ExecStart=/usr/bin/php artisan queue:work redis --sleep=3 --tries=3 --timeout=300 --max-time=3600
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
