# deploy/sportuniverse-scheduler.service
[Unit]
Description=Run SportsUniverse Laravel scheduler

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/sportuniverse
ExecStart=/usr/bin/php artisan schedule:run
